# php/supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=0.0.0.0:9001

[program:http_request_earliest]
command=php /var/www/html/public/index.php --command=AliReaza\\Atomic\\Commands\\HttpRequestEarliestListenerCommand
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
redirect_stderr=true
startsecs=0
exitcodes=0

[program:websocket_request_earliest]
command=php /var/www/html/public/index.php --command=AliReaza\\Atomic\\Commands\\WebSocketRequestEarliestListenerCommand
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
redirect_stderr=true
startsecs=0
exitcodes=0
